# SANsymphony/Proxmox Storage Plugin/SANsymphony-plugin_1.0.1_amd64/etc/udev/rules.d/99-datacore.rules
# Set timeout
SUBSYSTEM=="block", ACTION=="add|change", ATTRS{vendor}=="DataCore", ATTRS{model}=="Virtual Disk*", \
  RUN+="/bin/sh -c 'echo 80 > /sys/block/%k/device/timeout'"